# crontab/synchronized_lights
# Licensed under the BSD license.  See full license in LICENSE file.
# http://www.lightshowpi.org/
#
# Install these commands in the root user's crontab (i.e. sudo crontab -e)
#
# Make sure your system time zone is setup correctly. /var/log/syslog time
# stamps should match date command, if not reconfigure system time and reboot:
# sudo dpkg-reconfigure tzdata
# reboot

#set enviroment variable
# you might need to adjust the path if you did not install
# lightshowpi to the default directory
SYNCHRONIZED_LIGHTS_HOME=/home/pi/lightshowpi
﻿
# Start playing back songs and checking sms message at 6:30pm
30 18 * * * $SYNCHRONIZED_LIGHTS_HOME/bin/start_music_and_lights >> $SYNCHRONIZED_LIGHTS_HOME/logs/music_and_lights.play 2>&1 &

# Turn off the lights / sms at 10:15pm
15 22 * * * $SYNCHRONIZED_LIGHTS_HOME/bin/stop_music_and_lights >> $SYNCHRONIZED_LIGHTS_HOME/logs/music_and_lights.stop 2>&1 &
