# crontab/synchronized_lights
# Start playing back songs and checking sms message at 6:30pm
30 18 * * * exec /home/pi/bin/check_sms >> /var/log/music_and_lights.check 2>&1 &
30 18 * * * exec /home/pi/bin/play_sms >> /var/log/music_and_lights.play 2>&1 &

# Turn off the lights / sms at 10:15pm
15 22 * * * /home/pi/bin/stop_music_and_lights >> /var/log/music_and_lights.stop 2>&1 &


# BETA CMDS
exec $SYNCHRONIZED_LIGHTS_HOME/bin/check_sms >> $SYNCHRONIZED_LIGHTS_HOME/logs/music_and_lights.check 2>&1 &

exec $SYNCHRONIZED_LIGHTS_HOME/bin/play_sms >> $SYNCHRONIZED_LIGHTS_HOME/logs/music_and_lights.play 2>&1 &

$SYNCHRONIZED_LIGHTS_HOME/bin/stop_music_and_lights >> $SYNCHRONIZED_LIGHTS_HOME/logs/music_and_lights.stop 2>&1 &